# docker/php-cli/supervisor.conf
[program:pullTweets]
command = php %(ENV_APP_DIR)s/console.php twitter:test
numprocs = 1
autorestart = true
startsecs = 10
process_name = "%(program_name)s_%(process_num)s"
